# code-server-proxy/code-server-kill.service
[Unit]
Description=Watchdog for code-server

[Service]
Type=oneshot
ExecStart=/usr/local/bin/code-server-kill 2

