# code-server-proxy/code-server-kill.timer
[Unit]
Description=Watchdog for code-server

[Timer]
OnActiveSec=15m
WakeSystem=false
RemainAfterElapse=false

